# deploy/flashlink-webhook.service
[Unit]
Description=FlashLink VPN — YooKassa Webhook Server
After=network.target

[Service]
User=root
WorkingDirectory=/home/RemnaFlash
EnvironmentFile=/home/RemnaFlash/.env
ExecStart=/home/RemnaFlash/venv/bin/uvicorn webhook_server:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
